# install/projects_energy_monitor.service
[Unit]
 Description=Energy Monitor
 After=multi-user.target

[Service]
WorkingDirectory=/home/mnalavadi/energy_monitor
Type=idle
ExecStart=/home/mnalavadi/.local/bin/uv run app.py
User=mnalavadi

 [Install]
 WantedBy=multi-user.target
